# initramfs/var/spool/cron/crontabs/root
# Rotate logs daily
0 1 * * * /usr/sbin/logrotate /etc/logrotate.conf
